# config/wifi-notifier.service
[Unit]
Description=Aterm WiFi Client Notifier
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=your_user
Group=your_group
WorkingDirectory=/path/to/aterm-wifi-client-notifier
ExecStart=/usr/bin/python3 /path/to/aterm-wifi-client-notifier/src/wifi_notifier.py /path/to/aterm-wifi-client-notifier/config.json

# Restart on failure
Restart=on-failure
RestartSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/path/to/aterm-wifi-client-notifier

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wifi-notifier

[Install]
WantedBy=multi-user.target
